# contact-kiosk/deploy/contact-kiosk-https.service
[Unit]
Description=Contact directory kiosk backend (HTTPS, for camera access from other LAN machines)
After=network.target contact-kiosk.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/contact-kiosk
Environment=PATH=/home/pi/contact-kiosk/venv/bin
ExecStart=/home/pi/contact-kiosk/venv/bin/uvicorn app.main:app \
  --host 0.0.0.0 --port 8443 \
  --ssl-keyfile /home/pi/contact-kiosk/certs/kiosk.key \
  --ssl-certfile /home/pi/contact-kiosk/certs/kiosk.crt
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
